# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; ── PostgreSQL ───────────────────────────────────────────────
; Started first; the data dir is already initialised by entrypoint.sh.
[program:postgres]
command=/usr/local/bin/run-postgres.sh
user=postgres
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Spring Boot backend ──────────────────────────────────────
; The wrapper blocks until Postgres accepts connections, then exec's java,
; so Flyway/Hikari never see a cold database. DB host/user are fixed (the DB
; is internal to this container); DB_PASSWORD, OPENAI_API_KEY, AI_* and JWT_*
; are inherited from the container environment (set via .env / docker-compose),
; falling back to the defaults baked into application.yml when unset.
[program:backend]
command=/usr/local/bin/wait-for-postgres.sh java -jar /app/app.jar
directory=/app
user=root
priority=20
autostart=true
autorestart=true
startsecs=10
environment=DB_URL="jdbc:postgresql://127.0.0.1:5432/smart_quiz_hub",DB_USERNAME="postgres"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── nginx (serves the Angular build on :4200) ────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=30
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
